# src/usr/local/etc/rc.d/pondsec_ndr
#!/bin/sh

# PROVIDE: pondsec_ndr
# REQUIRE: LOGIN
# KEYWORD: shutdown

. /etc/rc.subr

name="pondsec_ndr"
rcvar="pondsec_ndr_enable"
command="/usr/sbin/daemon"
pidfile="/var/run/pondsec-ndr/pondsec-ndr.pid"
procname="/usr/local/bin/python3"
required_dirs="/usr/local/etc/pondsec-ndr /var/db/pondsec-ndr /var/log/pondsec-ndr /var/run/pondsec-ndr"

load_rc_config "$name"
: ${pondsec_ndr_enable:="NO"}
: ${pondsec_ndr_user:="pondsecndr"}
command_args="-p ${pidfile} -o /var/log/pondsec-ndr/pondsec-ndr-service.log /usr/local/sbin/pondsec-ndr"

run_rc_command "$1"
